# Agent/phi3-document-processor.service
[Unit]
Description=Phi-3 Document Processor Service
After=network.target ipfs.service
Wants=ipfs.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/757built.com
ExecStart=/bin/bash /root/757built.com/Agent/worker_startup.sh
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
